# nxlog/supervisord.conf
[supervisord]
nodaemon=true
 
[program:nxlog]
command=/usr/bin/nxlog -f
 
[program:logreader.py]
command=/opt/net2Host-dev/logreader.py winSecAudit nodaemon
