# nxlog/supervisord.conf
[supervisord]
nodaemon=true
 
[program:nxlog]
command=/usr/bin/nxlog -f
 
[program:winSecAudit]
command=/opt/net2Host/logreader.py winSecAudit nodaemon
